# scripts/com.seohyeong.lectureauto.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.seohyeong.lectureauto</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/seohyeongmin/Desktop/강의 홈페이지 제작/.venv/bin/python</string>
        <string>-u</string>
        <string>-m</string>
        <string>telegram_bot.poll</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/seohyeongmin/Desktop/강의 홈페이지 제작</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/lecture-daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/lecture-daemon.err.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>ko_KR.UTF-8</string>
    </dict>
</dict>
</plist>
